# systemd/setdns-route53.service
[Unit]
Description=sets route53 resource records based on config

[Service]
Type=oneshot
ExecStart=/usr/bin/setdns-route53.sh
